# systemd/loopover-docker-prune.service.example
# Runs scripts/selfhost-docker-prune.sh on a schedule (paired with loopover-docker-prune.timer.example) to
# reclaim unused Docker images and build cache before disk fills up on a long-running self-host deployment.
#
# Install (adjust the path below to wherever you cloned/deployed loopover):
#   sudo cp systemd/loopover-docker-prune.service.example /etc/systemd/system/loopover-docker-prune.service
#   sudo cp systemd/loopover-docker-prune.timer.example /etc/systemd/system/loopover-docker-prune.timer
#   sudo $EDITOR /etc/systemd/system/loopover-docker-prune.service   # fix WorkingDirectory / ExecStart path
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now loopover-docker-prune.timer
#
# This is a HOST-level unit, not a Docker Compose service: reclaiming images/build-cache needs real Docker
# daemon access, which no container in docker-compose.yml is granted (see that file's docker-proxy and
# runner service comments for why raw /var/run/docker.sock exposure into a container is avoided).

[Unit]
Description=LoopOver self-host Docker resource hygiene (image + build-cache prune)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# REQUIRED: point this at wherever you cloned loopover.
WorkingDirectory=/opt/loopover
ExecStart=/bin/sh /opt/loopover/scripts/selfhost-docker-prune.sh
# Optional: override the default 7-day (168h) safety window before an unused image/build-cache entry is
# eligible for removal.
# Environment=GITTENSORY_DOCKER_PRUNE_RETAIN_HOURS=168
